# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

# Run cron
/usr/sbin/cron -f -L 7
