# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/usr/bin/with-contenv bash

# Run cron
exec /usr/sbin/cron -f -L 7
